# system_files/usr/lib/systemd/system/systemd-suspend.service.d/disable_freeze_user_session.conf
# /usr/lib/systemd/system/systemd-suspend.service.d/disable_freeze_user_session.conf
#

[Service]
Environment="SYSTEMD_SLEEP_FREEZE_USER_SESSIONS=false"
